# php-base/supervisor/php-service.ini
[program:php]
command=/usr/local/bin/start-php-fpm.sh fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
